# systemd/aegis-fuse-daemon.service
[Unit]
Description=AEGIS FUSE Daemon (ADR-107)
Documentation=https://docs.100monkeys.ai/docs/architecture/storage-gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=RUST_LOG=info
Environment=AEGIS_FUSE_DAEMON_LISTEN_ADDR=0.0.0.0:50053
Environment=AEGIS_FUSE_DAEMON_ENDPOINT=http://host.containers.internal:50053
ExecStart=%h/100monkeys/aegis-deploy/bin/aegis fuse-daemon start \
  --orchestrator-url grpc://127.0.0.1:50051 \
  --mount-prefix /tmp/aegis-fuse-mounts \
  --listen-addr ${AEGIS_FUSE_DAEMON_LISTEN_ADDR}
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
